# php5-fpm.conf
[program:php5-fpm]
command=/usr/sbin/php5-fpm --nodaemonize --fpm-config /etc/php5/fpm/php-fpm.conf
process_name=%(program_name)s
numprocs=1
